# deploy/systemd/ata-city-agent.service
[Unit]
Description=ATA City Agent - Virtual City Simulator
After=network.target ata-soul-agent.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ata-city-agent --listen :8082 --db /var/lib/ata/city.db --souls-db /var/lib/ata/souls.db
Restart=on-failure
RestartSec=5
User=ata
Group=ata
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/ata
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ata-city-agent

[Install]
WantedBy=multi-user.target
